# doc_examples/rqworker-watcher.path
[Path]
PathModified=/home/mroa/SS2DTournamentRunner/main

[Install]
WantedBy=multi-user.target